# config/crontab
* * * * *   idetect     /home/idetect/scripts/update_mviews.bash
