# config/crontab
17 8 * * *   idetect     /home/idetect/scripts/update_mviews.bash
